# supervisor/celery.conf
[program:runbot]
process_name=%(program_name)s_%(process_num)02d
user=root
command=python3 /root/tgbot/manage.py run_bot
autostart=true
autorestart=true
redirect_stderr=true
